# check-plugins/service/unit-test/run
#!/usr/bin/env python3
# -*- coding: utf-8; py-indent-offset: 4 -*-
#
# Author:  Linuxfabrik GmbH, Zurich, Switzerland
# Contact: info (at) linuxfabrik (dot) ch
#          https://www.linuxfabrik.ch/
# License: The Unlicense, see LICENSE file.

# https://github.com/Linuxfabrik/monitoring-plugins/blob/main/CONTRIBUTING.md

import sys

sys.path.insert(0, '..')

import unittest

import lib.lftest
from lib.globals import STATE_OK, STATE_UNKNOWN, STATE_WARN

TESTS = [
    {
        'id': 'unknown-invalid-regex',
        'test': 'stdout/windows-services-mixed,,',
        'params': "--service '*'",
        'assert-retc': STATE_UNKNOWN,
        'assert-in': ['Invalid regex "*": nothing to repeat at position 0'],
    },
    {
        'id': 'ok-bfe-running',
        'test': 'stdout/windows-services-mixed,,0',
        'params': "--service='^bfe' --status=running --warning=1:1",
        'assert-retc': STATE_OK,
        'assert-in': [
            "Everything is ok. 1 service named r`^bfe` and start type ['automatic'] "
            "found, 1 in status ['running'] (thresholds 1:1/None).",
        ],
    },
    {
        'id': 'warn-win-services-stopped',
        'test': 'stdout/windows-services-mixed,,0',
        'params': (
            "--service='Win.*' --starttype=automatic --status=continue_pending "
            '--status=pause_pending --status=paused --status=start_pending '
            '--status=stop_pending --status=stopped --warning 0'
        ),
        'assert-retc': STATE_WARN,
        'assert-in': [
            "3 services named r`Win.*` and start type ['automatic'] found, "
            "1 in status ['continue_pending', 'pause_pending', 'paused', "
            "'start_pending', 'stop_pending', 'stopped'] (thresholds 0/None) [WARNING].",
        ],
    },

    # regression tests for issue #921: service names that contain a
    # space (e.g. "RAS Telegraf") must match just like any other name.
    # The fixture has the service with start_type=manual, so a default
    # call (--starttype=automatic implicit) must now give an
    # actionable error instead of the misleading "does not match any
    # service name".
    {
        'id': 'ok-ras-telegraf-manual',
        'test': 'stdout/windows-services-mixed,,0',
        'params': (
            "--service='^RAS Telegraf$' "
            '--starttype=manual --status=running --warning=1:1'
        ),
        'assert-retc': STATE_OK,
        'assert-in': [
            "Everything is ok. 1 service named r`^RAS Telegraf$` and start "
            "type ['manual'] found, 1 in status ['running'] "
            '(thresholds 1:1/None).',
            'RAS Telegraf',
        ],
    },
    {
        'id': 'unknown-ras-telegraf-default-starttype-filter',
        'test': 'stdout/windows-services-mixed,,0',
        # The user runs the plugin without thinking about --starttype;
        # the default filter [automatic] hides the RAS Telegraf service
        # (which is manual), but we now report that the regex did find
        # a candidate and tell them exactly how to include it.
        'params': "--service='^RAS Telegraf$'",
        'assert-retc': STATE_UNKNOWN,
        'assert-in': [
            'r`^RAS Telegraf$` matched 1 service',
            "start_type in ['automatic']",
            "Matched start_types: ['manual']",
            'RAS Telegraf (manual)',
            'Adjust --starttype',
        ],
    },
    {
        'id': 'unknown-regex-does-not-match-any-service',
        'test': 'stdout/windows-services-mixed,,0',
        # The regex really does not match any service in the fixture.
        # The old "does not match any service name" wording is kept for
        # this case so existing monitoring setups that grep on the
        # error text still work.
        'params': "--service='^doesnotexist$'",
        'assert-retc': STATE_UNKNOWN,
        'assert-in': [
            'r`^doesnotexist$` does not match any service name.',
        ],
    },
]


class TestCheck(unittest.TestCase):

    check = '../service'


lib.lftest.attach_tests(TestCheck, TESTS)


if __name__ == '__main__':
    unittest.main()
